# supervisor.ini
[supervisord]
nodaemon = true

[program:wsgi]
user = app
command = /home/app/start-gunicorn.sh
environment=PATH="/srv/virtenv/bin:%(ENV_PATH)s"
autorestart = false
startretries=0
stdout_events_enabled = true
stderr_events_enabled = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command = /usr/bin/start-nginx.sh
autorestart = false
startretries=0
stdout_events_enabled = true
stderr_events_enabled = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler

[eventlistener:fatal]
command = "supervisor-shutdown"
events = PROCESS_STATE_FATAL

